# files/upstart.conf
# launch freepuppet-run on host startup
start on startup
exec /usr/local/bin/freepuppet-run